# layers/meta-hl/recipes-support/hl-rsm/hl-rsm/hl-rsm.timer
[Unit]
Description=RSM client check timer

[Timer]
OnBootSec=1min
OnUnitInactiveSec=7min

[Install]
WantedBy=resin.target
